# root/etc/services.d/webdav/run
#!/usr/bin/with-contenv bash

if [[ -n "$WEBDAV_URL" ]]; then
    echo "**** Starting WebDAV mount service with rclone ****"
    
    # 创建挂载点
    mkdir -p "$WEBDAV_MOUNT_PATH"
    
    # 创建 rclone 配置
    mkdir -p /config/rclone
    
    cat > /config/rclone/rclone.conf << EOF
[webdav]
type = webdav
url = $WEBDAV_URL
vendor = other
user = $WEBDAV_USERNAME
pass = $(echo -n "$WEBDAV_PASSWORD" | rclone obscure -)
EOF
    
    # 等待网络就绪
    sleep 5
    
    # 挂载 WebDAV
    echo "**** Mounting WebDAV with rclone ****"
    exec rclone mount webdav: "$WEBDAV_MOUNT_PATH" \
        --config /config/rclone/rclone.conf \
        --allow-other \
        --allow-non-empty \
        --vfs-cache-mode writes \
        --log-level INFO
else
    echo "**** No WebDAV URL provided, WebDAV service disabled ****"
    sleep infinity
fi
